# modules/init.d/nix-daemon
#!/sbin/openrc-run

description="Nix multi-user build daemon"

command="/nix/var/nix/profiles/default/bin/nix-daemon"
command_background="yes"
pidfile="/run/nix-daemon.pid"

depend() {
    # Use 'after localmount' instead of 'need' to avoid hard-fail on some Devuan setups
    after localmount
    after bootmisc
    use net
}

start_pre() {
    # Ensure the socket directory exists with correct permissions
    checkpath -d -m 0755 /nix/var/nix/daemon-socket
}
