# modules/init.d/nix-daemon
#!/sbin/openrc-run

description="Nix multi-user build daemon"

command="/nix/var/nix/profiles/default/bin/nix-daemon"
command_background="yes"
pidfile="/run/nix-daemon.pid"

depend() {
    # 'use net' is less strict than 'need net'
    # It will start even if the virtual 'net' service isn't explicitly provided
    use net
    after network
    before dbus
}

start_pre() {
    checkpath -d -m 0755 /nix/var/nix/daemon-socket
}
