# modules/init.d/nix-daemon
#!/sbin/openrc-run

description="Nix multi-user build daemon"

command="/nix/var/nix/profiles/default/bin/nix-daemon"
command_background="yes"
pidfile="/run/nix-daemon.pid"

depend() {
    # Nix needs the filesystem to be mounted
    need localmount
    # We want it after basic system initialization
    after bootmisc
    # Network is optional but useful
    use net
}

start_pre() {
    # Ensure the socket directory exists with correct permissions
    checkpath -d -m 0755 /nix/var/nix/daemon-socket
}
